# launchd/com.finsuccess.stockmonitor.news-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.finsuccess.stockmonitor.news-update</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Frameworks/Python.framework/Versions/3.11/bin/python3</string>
		<string>/Users/thanhnhan/Documents/FinSuccess/Stock monitor/auto_update.py</string>
		<string>news</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/tmp</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>1800</integer>
	<key>StandardOutPath</key>
	<string>/Users/thanhnhan/Documents/FinSuccess/Stock monitor/logs/launchd-news-update.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/thanhnhan/Documents/FinSuccess/Stock monitor/logs/launchd-news-update.err.log</string>
</dict>
</plist>
